# assets/game/scripts/cabal_services/cabal_server1.conf
[program:DBAgent_01]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/DBAgent_01 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=3


[program:ChatDBAgent_01]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/ChatDBAgent_01 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=3


[program:CashDBAgent_01]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/CashDBAgent_01 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=3

[program:PartySvr_01]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/PartySvr_01 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=3

[program:ChatNode_01]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/ChatNode_01 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=3

[program:AgentShop_01]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/AgentShop_01 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=3






[program:WorldSvr_01_01]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'LD_PRELOAD="/usr/lib/cabal/lib.so /usr/lib/cabal/libchilkat-9.5.0.so /usr/lib/cabal/libPolyHook_2.so /usr/lib/cabal/libZycore.so.1.4 /usr/lib/cabal/libZydis.so.4.0" exec /usr/bin/WorldSvr_01_01 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_02]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'LD_PRELOAD="/usr/lib/cabal/lib.so /usr/lib/cabal/libchilkat-9.5.0.so /usr/lib/cabal/libPolyHook_2.so /usr/lib/cabal/libZycore.so.1.4 /usr/lib/cabal/libZydis.so.4.0" exec /usr/bin/WorldSvr_01_02 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_03]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'LD_PRELOAD="/usr/lib/cabal/lib.so /usr/lib/cabal/libchilkat-9.5.0.so /usr/lib/cabal/libPolyHook_2.so /usr/lib/cabal/libZycore.so.1.4 /usr/lib/cabal/libZydis.so.4.0" exec /usr/bin/WorldSvr_01_03 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_04]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'LD_PRELOAD="/usr/lib/cabal/lib.so /usr/lib/cabal/libchilkat-9.5.0.so /usr/lib/cabal/libPolyHook_2.so /usr/lib/cabal/libZycore.so.1.4 /usr/lib/cabal/libZydis.so.4.0" exec /usr/bin/WorldSvr_01_04 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_05]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'LD_PRELOAD="/usr/lib/cabal/lib.so /usr/lib/cabal/libchilkat-9.5.0.so /usr/lib/cabal/libPolyHook_2.so /usr/lib/cabal/libZycore.so.1.4 /usr/lib/cabal/libZydis.so.4.0" exec /usr/bin/WorldSvr_01_05 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_06]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'LD_PRELOAD="/usr/lib/cabal/lib.so /usr/lib/cabal/libchilkat-9.5.0.so /usr/lib/cabal/libPolyHook_2.so /usr/lib/cabal/libZycore.so.1.4 /usr/lib/cabal/libZydis.so.4.0" exec /usr/bin/WorldSvr_01_06 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_07]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'LD_PRELOAD="/usr/lib/cabal/lib.so /usr/lib/cabal/libchilkat-9.5.0.so /usr/lib/cabal/libPolyHook_2.so /usr/lib/cabal/libZycore.so.1.4 /usr/lib/cabal/libZydis.so.4.0" exec /usr/bin/WorldSvr_01_07 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_08]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'LD_PRELOAD="/usr/lib/cabal/lib.so /usr/lib/cabal/libchilkat-9.5.0.so /usr/lib/cabal/libPolyHook_2.so /usr/lib/cabal/libZycore.so.1.4 /usr/lib/cabal/libZydis.so.4.0" exec /usr/bin/WorldSvr_01_08 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_09]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'LD_PRELOAD="/usr/lib/cabal/lib.so /usr/lib/cabal/libchilkat-9.5.0.so /usr/lib/cabal/libPolyHook_2.so /usr/lib/cabal/libZycore.so.1.4 /usr/lib/cabal/libZydis.so.4.0" exec /usr/bin/WorldSvr_01_09 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_10]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'LD_PRELOAD="/usr/lib/cabal/lib.so /usr/lib/cabal/libchilkat-9.5.0.so /usr/lib/cabal/libPolyHook_2.so /usr/lib/cabal/libZycore.so.1.4 /usr/lib/cabal/libZydis.so.4.0" exec /usr/bin/WorldSvr_01_10 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_11]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'LD_PRELOAD="/usr/lib/cabal/lib.so /usr/lib/cabal/libchilkat-9.5.0.so /usr/lib/cabal/libPolyHook_2.so /usr/lib/cabal/libZycore.so.1.4 /usr/lib/cabal/libZydis.so.4.0" exec /usr/bin/WorldSvr_01_11 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_12]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'LD_PRELOAD="/usr/lib/cabal/lib.so /usr/lib/cabal/libchilkat-9.5.0.so /usr/lib/cabal/libPolyHook_2.so /usr/lib/cabal/libZycore.so.1.4 /usr/lib/cabal/libZydis.so.4.0" exec /usr/bin/WorldSvr_01_12 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_13]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'LD_PRELOAD="/usr/lib/cabal/lib.so /usr/lib/cabal/libchilkat-9.5.0.so /usr/lib/cabal/libPolyHook_2.so /usr/lib/cabal/libZycore.so.1.4 /usr/lib/cabal/libZydis.so.4.0" exec /usr/bin/WorldSvr_01_13 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_14]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'LD_PRELOAD="/usr/lib/cabal/lib.so /usr/lib/cabal/libchilkat-9.5.0.so /usr/lib/cabal/libPolyHook_2.so /usr/lib/cabal/libZycore.so.1.4 /usr/lib/cabal/libZydis.so.4.0" exec /usr/bin/WorldSvr_01_14 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_15]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'LD_PRELOAD="/usr/lib/cabal/lib.so /usr/lib/cabal/libchilkat-9.5.0.so /usr/lib/cabal/libPolyHook_2.so /usr/lib/cabal/libZycore.so.1.4 /usr/lib/cabal/libZydis.so.4.0" exec /usr/bin/WorldSvr_01_15 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_16]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'LD_PRELOAD="/usr/lib/cabal/lib.so /usr/lib/cabal/libchilkat-9.5.0.so /usr/lib/cabal/libPolyHook_2.so /usr/lib/cabal/libZycore.so.1.4 /usr/lib/cabal/libZydis.so.4.0" exec /usr/bin/WorldSvr_01_16 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_17]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'LD_PRELOAD="/usr/lib/cabal/lib.so /usr/lib/cabal/libchilkat-9.5.0.so /usr/lib/cabal/libPolyHook_2.so /usr/lib/cabal/libZycore.so.1.4 /usr/lib/cabal/libZydis.so.4.0" exec /usr/bin/WorldSvr_01_17 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_18]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'LD_PRELOAD="/usr/lib/cabal/lib.so /usr/lib/cabal/libchilkat-9.5.0.so /usr/lib/cabal/libPolyHook_2.so /usr/lib/cabal/libZycore.so.1.4 /usr/lib/cabal/libZydis.so.4.0" exec /usr/bin/WorldSvr_01_18 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_19]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'LD_PRELOAD="/usr/lib/cabal/lib.so /usr/lib/cabal/libchilkat-9.5.0.so /usr/lib/cabal/libPolyHook_2.so /usr/lib/cabal/libZycore.so.1.4 /usr/lib/cabal/libZydis.so.4.0" exec /usr/bin/WorldSvr_01_19 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_20]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'LD_PRELOAD="/usr/lib/cabal/lib.so /usr/lib/cabal/libchilkat-9.5.0.so /usr/lib/cabal/libPolyHook_2.so /usr/lib/cabal/libZycore.so.1.4 /usr/lib/cabal/libZydis.so.4.0" exec /usr/bin/WorldSvr_01_20 -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15
